# autostart/autostart.sh.desktop
[Desktop Entry]
Exec=~/bin/autostart.sh
Icon=
Name=autostart.sh
Path=
Terminal=True
Type=Application
